# deploy/maintenance/docker-prune.service
[Unit]
Description=Prune unused Docker images, containers, and BuildKit cache
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-prune.sh
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
